# etc/supervisord.d/hbase.ini
[program:hbase]
command = /opt/qnib/hbase/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
