# etc/supervisord.d/hbase.ini
[program:hbase]
command = /opt/qnib/hbase/bin/start_hbase.sh
redirect_stderr=true
stdout_syslog=true
